# master-fs/lib/udev/rules.d/60-master-usb-if.rules
SUBSYSTEM=="net", RUN+="master-to-terminal-if.agent"
